# packages/runit/srv-wg/services/wg/run
#!/bin/sh

exec 2>&1

[ -r ./conf ] && . ./conf

echo 'wg starting.'

LANG=C LC_ALL=C \
/usr/bin/wg-quick up $INTERFACE 2>&1

exec pause