# packages/systemd/edgevpn-systemd/edgevpn.service
[Unit]
Description=EdgeVPN Daemon
After=network.target

[Service]
ExecStart=/usr/bin/edgevpn
Restart=always

[Install]
WantedBy=multi-user.target
