# packages/systemd/edgevpn-systemd/edgevpn.service
[Unit]
Description=EdgeVPN Daemon
After=network.target

[Service]
EnvironmentFile=/etc/systemd/system.conf.d/edgevpn.env
ExecStart=/usr/bin/edgevpn
Restart=always

[Install]
WantedBy=multi-user.target
