# packages/runit/srv-sshd/services/sshd/run
#!/bin/sh

exec 2>&1

[ -r ./conf ] && . ./conf

echo 'sshd starting.'

LANG=C LC_ALL=C \
exec /usr/sbin/sshd -D 2>&1
